# example.mixin/etc/systemd/system/daytime@.service
[Unit]
Description=Daytime service

[Service]
ExecStart=/bin/date
StandardInput=socket
StandardOutput=socket
