# test.mixin/etc/systemd/system/test-no-pass-fds.service
[Unit]
Description=Test service

[Service]
ExecStart=/test/test_pass_listen_fds.py expect four
Type=oneshot
